# ACCREDiagnostics/testcron.crontab
* * * * * /home/janjamrk/ACCREDiagnostics/simple.sh
